# packages/stfu/stfu.install
post_install() {
    cat <<EOF

To enable the STFU daemon for you user run the following command:
systemctl --user enable stfu.service --now
EOF
}

