# packages/stfu/stfu.install
post_install() {
    cat <<EOF

Enable stfu in TTY-Only mode
systemctl --user enable stfu.service --now

Enable stfu for a desktop
ln -sf /usr/share/stfu/stfu.desktop ~/.config/autostart/stfu.desktop
EOF
}

